# systemd/emacsd.service
[Unit]
Description=emacs server daemon

[Service]
ExecStart=/usr/bin/emacs --daemon

[Install]
WantedBy=daemon.target
