# systemd/urxvtd.service
[Unit]
Description=urxvt daemon

[Service]
ExecStart=/usr/bin/urxvtd --opendisplay

[Install]
WantedBy=daemon.target
