# systemd/X11.target
[Unit]
Description=graphical X11 session
