# systemd/user/daemon.target
[Unit]
Description=general user daemons

[Install]
WantedBy=graphical-session.target
