# systemd/user/daemon.target
[Unit]
Description=general user daemons
AllowIsolate=yes
Wants=default.target
