# systemd/user/daemon.target
[Unit]
Description=user daemons

[Install]
Alias=default.target

