# systemd/user/dunst.service
[Unit]
Description=graphical notification daemon

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=X11.target
