# systemd/user/wayland.target
[Unit]
Description=wayland-based session daemons
