# systemd/user/X11.target
[Unit]
Description=X11-based session daemons
