# tty/loadkeys.service
[Unit]
Description=Load TTY keymap tweaks

[Service]
Type=simple
ExecStart=/usr/bin/loadkeys /etc/vkeymap.conf

[Install]
WantedBy=getty.target
