# tty/setvtrgb.service
[Unit]
Description=Set TTY color scheme

[Service]
Type=simple
ExecStart=/usr/bin/setvtrgb /etc/vcolormap.conf

[Install]
WantedBy=getty.target
